# templates/default/devpi.init.erb
#!/bin/sh
# Devpi server startup script 
#
# DO NOT MODIFY THIS FILE
# Dynamically generated by Chef on <%= node['fqdn'] %>
# Local modifications will be overwritten by Chef.
#
#chkconfig: 2345 80 05
#description: Atlassian bamboo

### BEGIN INIT INFO
# Provides:          devpi
# Required-Start:    $local_fs $remote_fs $network $syslog $named
# Required-Stop:     $local_fs $remote_fs $network $syslog $named
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# X-Interactive:     true
# Short-Description: Start/stop devpi server
### END INIT INFO

# Source function library.
<% case node['platform_family'] -%>
<% when 'rhel' -%>
. /etc/rc.d/init.d/functions
<% else -%>
. /lib/lsb/init-functions
<% end -%>

APP=<%= @name %>
USER=<%= @user %>

BASE=<%= @home_dir %>
DATA=<%= @data_dir %>

HOST=<%= @host %>
PORT=<%= @port %>

CMD="$BASE/bin/python $BASE/bin/$APP --serverdir $DATA"


start() {
  if test ! -e $DATA/.nodeinfo; then
    echo "Initializating $APP"
    su -l $USER -c "$CMD --init"
  fi
  echo "Starting $APP"
  su -l $USER -c "$CMD --host $HOST --port $PORT --start"
  return $?
}

stop() {
  echo "Stopping $APP"
  su -l $USER -c "$CMD --stop"
  return $?
}

case "$1" in
  # Start command
  start)
    start
    RETVAL=$?
    ;;
  # Stop command
  stop)
    stop
    RETVAL=$?
    ;;
  # Restart command
  restart)
    stop
    start
    RETVAL=$?
    ;;
  # Status command
  status)
    $CMD --status
    RETVAL=$?
    ;;
  *)
    echo "Usage: $APP {start|restart|stop|status}"
    RETVAL=2
    ;;
esac

exit $RETVAL
