# templates/default/devpi.service.erb
[Unit]
Description=Devpi Server
After=syslog.target network.target

[Service]
Type=simple
User=<%= @user %>
Group=<%= @group %>
ExecStart=source <%= @home_dir %>/bin/activate && <%= @home_dir %>/bin/<%= @name %> --host <%= @url %> --port <%= @port %> --serverdir <%= @data_dir %>

[Install]
WantedBy=multi-user.target
